# home/dot_config/systemd/user/sync-stg.service
[Unit]
Description=Syncs tabs of Firefox addon STG in gist

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/true
ExecStop=%h/scripts/sync-stg.sh

[Install]
WantedBy=multi-user.target